# linbofs/etc/inittab
# /etc/inittab
#
# thomas@linuxmuster.net
# 20231207
#
::sysinit:/init.sh
::respawn:/linbo.sh
::wait:/usr/bin/linbo_vnc onboot
